# crontab
* * * * * bundle exec rake heroku:test_cron