# crontab
0 0 * * 0 bundle exec rake vita_providers:scrape_vita_providers