# crontab
0 0 * * 0 bundle exec rake vita_providers:scrape_vita_providers
* * * * * bundle exec rake search:refresh
* * * * * bundle exec rake stats:track_metrics
20 21 * * 1-5 bundle exec rake client_surveys:send_client_in_progress_surveys
15 20 * * * bundle exec rake not_ready:remind
*/10 * * * * bundle exec rake efile:poll_and_get_acknowledgments
