# crontab
*/30 * * * * bundle exec rake documents:notify_unsent